# systemd/auth-fw-passwd-set.socket.in
[Socket]
ListenStream=@RUN_DIR@/.@SOCK_PASSWD_SET@
SocketUser=security_fw
SocketGroup=security_fw
SocketMode=0777
SmackLabelIPIn=*
SmackLabelIPOut=@

Service=@SERVICE_NAME@.service

[Unit]
Wants=@SERVICE_NAME@.target
Before=@SERVICE_NAME@.target

[Install]
WantedBy=sockets.target
