# systemd/auth-fw.service.in
[Unit]
Description=Start @SERVICE_NAME@ server

[Service]
Type=notify
ExecStart=@BIN_DIR@/@SERVICE_NAME@
User=security_fw
Group=security_fw
Sockets=@SOCK_PASSWD_CHECK@
Sockets=@SOCK_PASSWD_SET@
Sockets=@SOCK_PASSWD_RESET@
Sockets=@SOCK_PASSWD_POLICY@

[Install]
WantedBy=multi-user.target
