# systemd/authentication-server.target
[Unit]
Description=authentication server sockets
DefaultDependencies=true

